# musicbot.service
[Unit]
Description=MusicBot TypeScript Runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=musicbot
Group=musicbot
WorkingDirectory=/opt/musicbot/ts-bot
EnvironmentFile=/etc/musicbot/musicbot.env
ExecStart=/usr/bin/node /opt/musicbot/ts-bot/dist/main.js
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
